# scripts/deploy/wecom-proxy.service
[Unit]
Description=WeCom Relay Proxy - 企业微信消息中转代理
Documentation=https://github.com/SuLvXiangXin/arxiv_daily
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/wecom-proxy
ExecStart=/usr/bin/node /opt/wecom-proxy/wecom_proxy.js
Restart=always
RestartSec=5
Environment=PORT=9000
EnvironmentFile=/opt/wecom-proxy/.env

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wecom-proxy

# 安全
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/wecom-proxy

[Install]
WantedBy=multi-user.target
